# provisioning/roles/consul/templates/etc/profile.d/consul.sh.j2
export PATH={{ consul_binary_dir }}:${PATH}
export CONSUL_HOME={{ consul_home }}
